# posa/src/deb/control/postrm
#!/bin/sh -e
# postrm for posa

if [ "$1" = "purge" ] ; then
    update-rc.d posa remove >/dev/null
fi

exit 0
